# unitfiles/aurbuilder-package.service
[Unit]
Description=aurbuilder Package build job
OnSuccess=aurbuilder-repo.service

[Service]
Type=simple
User=aurbuilder
Group=aurbuilder
ExecStart=/usr/lib/aurbuilder/build_packages.sh
Restart=no
StartLimitInterval=10min
StartLimitBurst=8
