# unitfiles/aurbuilder-package.timer
[Unit]
Description=Daily AUR package rebuild

[Timer]
OnCalendar=4:00
Persistent=true

[Install]
WantedBy=timers.target
