# unitfiles/aurbuilder-container.timer
[Unit]
Description=Time for aurbuilder container rebuild

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
